# cookbooks/quantum/templates/default/quantum-metadata-agent.conf.erb
description "Quantum metadata agent"
author "whitekid@gmail.com"

start on runlevel [2345]
stop on runlevel [016]

chdir /var/run

pre-start script
	mkdir -p /var/run/quantum
	chown quantum:root /var/run/quantum
end script

exec start-stop-daemon --start --chuid quantum --exec /usr/local/bin/quantum-metadata-agent -- --config-file=/etc/quantum/quantum.conf --config-file=/etc/quantum/metadata_agent.ini --log-file=/var/log/quantum/metadata-agent.log
