# cookbooks/repo/templates/default/ubuarchive.cron.erb
#
# Regular cron jobs for the ubuarchive package
#
0 4     * * *   root      ubuarchive

